# fs/etc/supervisor/conf.d/docker.conf
[program:docker]
priority=10
command=dind wrapdocker
startsecs=0
exitcodes=0,1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0